# com.davidmackinnon.brief-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd job: back up the brief clone to your GitHub each day at 18:15 local. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.davidmackinnon.brief-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/davidmackinnon/tradingview-mcp-jackson/autocommit-brief.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
